# supervisord.conf
[supervisord]
nodaemon=true

[supervisorctl]

[program:django]
command=python3 /code/manage.py runserver 0.0.0.0:8000

[program:tgbot]
command=python3 /code/manage.py bot